# recipes-core/systemd-conf/files/bootup_init.service
[Unit]
Description= sima bootup init 

[Service]
Type=oneshot
ExecStart=/usr/bin/bootup_init.sh

[Install]
WantedBy=multi-user.target

